# systemd/dfsort.service
[Unit]
Description=DFSORT File Auto-Sorter
After=network.target

[Service]
Type=simple
User=%i
ExecStart=/usr/bin/dfsort --daemon
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target